# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:php-fpm]
command=/usr/sbin/php-fpm8.1 -F
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:php-error-log]
command=/usr/bin/tail -f /var/www/meshsilo/storage/logs/php-error.log
autostart=true
autorestart=true
priority=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:app-log]
command=/usr/bin/tail -f /var/www/meshsilo/storage/logs/app.log
autostart=true
autorestart=true
priority=16
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:security-log]
command=/usr/bin/tail -f /var/www/meshsilo/storage/logs/security.log
autostart=true
autorestart=true
priority=17
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:access-log]
command=/usr/bin/tail -f /var/www/meshsilo/storage/logs/access.log
autostart=true
autorestart=true
priority=18
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Per-queue workers — each one is single-threaded, but they run in parallel
# with each other. A long STL->3MF conversion in the `conversions` worker
# does NOT block a quick thumbnail in the `thumbnails` worker.
#
# Memory footprint: each worker is a full PHP CLI process (~60-80MB idle,
# bigger during heavy jobs). Six workers is ~400MB idle — budget accordingly.
# Reduce by removing queues you don't use.

# ProcessUpload — zip extraction after tus completion. Disk/I-O heavy.
[program:queue-worker-uploads]
command=/usr/bin/php /var/www/meshsilo/cli/queue-worker.php --queue=uploads --sleep=3 --max-time=3600 --max-jobs=200
user=www-data
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ConvertStlTo3mf — memory-heavy. Runs 3 workers in parallel so multiple
# STL conversions can proceed at once. Queue::pop() uses a DB transaction
# to reserve jobs atomically, so the 3 workers won't double-dispatch the
# same job.
#
# Memory safety: each worker runs the preflight check in
# convertPartTo3MF() (includes/converter.php) that looks at cgroup memory
# and throws if usage is >= 90%. When RAM is tight, extra workers will
# back off and retry rather than fight for allocations. If you see
# frequent 'System memory too high' failures, reduce numprocs back to 1
# or bump the container mem_limit in docker-compose.yml.
[program:queue-worker-conversions]
command=/usr/bin/php /var/www/meshsilo/cli/queue-worker.php --queue=conversions --sleep=5 --max-time=3600 --max-jobs=100
process_name=%(program_name)s_%(process_num)02d
numprocs=3
user=www-data
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# GenerateThumbnail — 3mf embedded-thumbnail extraction + OpenSCAD renders.
[program:queue-worker-thumbnails]
command=/usr/bin/php /var/www/meshsilo/cli/queue-worker.php --queue=thumbnails --sleep=5 --max-time=3600 --max-jobs=200
user=www-data
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# OptimizeImage — JPEG/PNG → WebP. Light.
[program:queue-worker-images]
command=/usr/bin/php /var/www/meshsilo/cli/queue-worker.php --queue=images --sleep=5 --max-time=3600 --max-jobs=500
user=www-data
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# OptimizePdf — Ghostscript or qpdf. Medium weight.
[program:queue-worker-pdfs]
command=/usr/bin/php /var/www/meshsilo/cli/queue-worker.php --queue=pdfs --sleep=5 --max-time=3600 --max-jobs=200
user=www-data
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# default — catches anything not routed to a specific queue (SendEmailJob,
# ad-hoc jobs, plugin-dispatched jobs that don't specify a queue).
[program:queue-worker-default]
command=/usr/bin/php /var/www/meshsilo/cli/queue-worker.php --queue=default --sleep=5 --max-time=3600 --max-jobs=200
user=www-data
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cron]
command=/bin/sh -c "while true; do /usr/bin/php /var/www/meshsilo/cli/cron.php; sleep 60; done"
user=www-data
autostart=true
autorestart=true
priority=25
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
